# examples/breezyd.service
# Systemd unit for breezyd, for Linux distros that don't use NixOS.
#
# This mirrors the unit the NixOS module generates from nix/module.nix —
# the Nix attrset there is the source of truth for sandboxing tweaks
# (e.g. RestrictAddressFamilies needs AF_NETLINK for HomeKit mDNS to
# work). When in doubt, cross-reference nix/module.nix.
#
# Install path: /etc/systemd/system/breezyd.service
# Then: sudo systemctl daemon-reload && sudo systemctl enable --now breezyd
#
# See README.md "Linux + systemd" for the surrounding setup (user
# creation, config placement, firewall).

[Unit]
Description=Vents Twinfresh Breezy / Elite 160 Pro daemon
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/local/bin/breezyd --config /etc/breezyd/breezyd.toml
User=breezyd
Group=breezyd
Restart=on-failure
RestartSec=5s

# Energy state files (and HomeKit pairing files, when [homekit].enabled
# is true) live under /var/lib/breezyd/. systemd creates and chowns it
# to the service user. Safe to leave enabled even without HomeKit.
StateDirectory=breezyd

# Hardening — breezyd only needs outbound UDP and an HTTP listener.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
# AF_NETLINK is needed by Go's net.Interfaces() on Linux. Without it
# the HomeKit bridge runs but advertises on zero interfaces (mDNS-deaf);
# iPhones never see the bridge. No log, no error; just silence on
# UDP/5353.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged

[Install]
WantedBy=multi-user.target
